# services_example/chat_backend.service
[Unit]
Description=Flask Chat Backend
After=network.target

[Service]
Type=simple
User=agente
WorkingDirectory=/opt/agente/backend
ExecStart=/opt/agente/venv/bin/python chat_backend.py
Restart=always
Environment="PATH=/opt/agente/venv/bin"

[Install]
WantedBy=multi-user.target

