# services_example/embed_backend.service
[Unit]
Description=Flask Embeddings Backend
After=network.target

[Service]
Type=simple
User=agente
WorkingDirectory=/opt/agente/backend
ExecStart=/opt/agente/venv/bin/python embed_backend.py
Restart=always
Environment="PATH=/opt/agente/venv/bin"

[Install]
WantedBy=multi-user.target

